# recipes-support/gpio-service/files/gpio.service
[Service]
ExecStart=/usr/sbin/gpio.sh
Restart=always
[Install]
WantedBy=multi-user.target
[Service]
StandardOutput=null
StandardError=null